# systemd/audiobooks-proxy.service
[Unit]
Description=Audiobooks Library HTTPS Reverse Proxy
Documentation=https://github.com/greogory/Audiobook-Manager
After=network-online.target audiobooks-api.service
Requires=audiobooks-api.service
Wants=network-online.target

[Service]
Type=simple
User=audiobooks
Group=audiobooks
WorkingDirectory=/opt/audiobooks/library/web-v2

# Load configuration from /etc/audiobooks/audiobooks.conf
EnvironmentFile=/etc/audiobooks/audiobooks.conf
Environment=PYTHONUNBUFFERED=1

# Only start if port is not already in use
ExecStartPre=/bin/sh -c '! /usr/bin/lsof -i:${AUDIOBOOKS_WEB_PORT} >/dev/null 2>&1'
ExecStart=/usr/bin/python3 /opt/audiobooks/library/web-v2/proxy_server.py

Restart=on-failure
RestartSec=5
# Use journal logging (query with: journalctl -u audiobooks-proxy)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audiobooks-proxy

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadOnlyPaths=/etc/audiobooks /opt/audiobooks
ReadWritePaths=/var/lib/audiobooks  # default
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=audiobooks.target multi-user.target
