# systemd/audiobooks.target
[Unit]
Description=Audiobook Library Services
After=network.target

[Install]
WantedBy=multi-user.target
