# systemd/audiobooks-downloader.service
[Unit]
Description=Audiobook Downloader (Audible sync)
After=network-online.target
Wants=network-online.target
PartOf=audiobooks.target

[Service]
Type=oneshot
User=audiobooks
Group=audiobooks

# Load configuration
EnvironmentFile=-/etc/audiobooks/audiobooks.conf
Environment=HOME=/var/lib/audiobooks

# Create runtime directories
RuntimeDirectory=audiobooks
RuntimeDirectoryMode=0755

# Main script
ExecStart=/usr/local/bin/audiobooks-download

# Don't restart automatically (triggered by timer)
Restart=no

[Install]
WantedBy=audiobooks.target
