# systemd/audiobooks-periodicals-sync.timer
[Unit]
Description=Sync Audible periodicals index twice daily
Documentation=file:///opt/audiobooks/docs/PERIODICALS.md

[Timer]
# Run at 6 AM and 6 PM daily
OnCalendar=*-*-* 06:00:00
OnCalendar=*-*-* 18:00:00

# Catch up if system was off during scheduled time
Persistent=true

# Add random delay up to 5 minutes to avoid API hammering
RandomizedDelaySec=300

# Don't wake system from sleep just for this
WakeSystem=false

[Install]
WantedBy=audiobooks.target timers.target
