# systemd/audiobooks-api.service
[Unit]
Description=Audiobooks Library API Server (Waitress)
Documentation=https://github.com/greogory/Audiobook-Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=audiobooks
Group=audiobooks
WorkingDirectory=/opt/audiobooks/library/backend

# Load configuration from /etc/audiobooks/audiobooks.conf
EnvironmentFile=/etc/audiobooks/audiobooks.conf
Environment=PYTHONUNBUFFERED=1

# Only start if port is not already in use
ExecStartPre=/bin/sh -c '! /usr/bin/lsof -i:${AUDIOBOOKS_API_PORT} >/dev/null 2>&1'
ExecStart=/opt/audiobooks/library/venv/bin/python api_server.py

Restart=on-failure
RestartSec=5
# Use journal logging for better log management (query with: journalctl -u audiobooks-api)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audiobooks-api

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/audiobooks  # default
ProtectHome=yes
# PrivateTmp=no - API needs access to /tmp/audiobook-staging (shared with converter)
PrivateTmp=no

[Install]
WantedBy=audiobooks.target multi-user.target
